# systemd/telegram-inbox-bot.service
[Unit]
Description=Telegram Inbox Bot Telegram Polling Service
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=telegraminbox
Group=telegraminbox
WorkingDirectory=/var/www/contactinboxbot
EnvironmentFile=/var/www/contactinboxbot/.env
ExecStart=/var/www/contactinboxbot/.venv/bin/python -m app.bot.main
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/contactinboxbot

[Install]
WantedBy=multi-user.target
